# systemd_units/cg_scraper.timer
[Unit]
Description=Ejecutar script de RAM 3 veces al día

[Timer]
OnCalendar=*-*-* 15,17,19:00:00
Persistent=true

[Install]
WantedBy=timers.target
